# Arhangelskiy teatr/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::restart:/sbin/init

::once:/controller/protocolvm
::once:/usr/sbin/lighttpd -D -f /etc/lighttpd.conf
::once:/gui/run
::shutdown:/bin/umount -a -r
